# devops/crontab
* * * * *  php /home/jtodd/intrusionDetect/artisan util:monitor
